# automated-commands.cron
@reboot ps aux | grep volume_controller_evdev_lirc | grep -v grep | awk '{print $2}' | xargs kill; nohup python3.11 ~/code/volume-control/scripts/volume_controller_evdev_lirc.py &> /tmp/nohup.out & disown
