# automated-commands.cron
# this is currently not utilized
@reboot ps aux | grep volume_control | grep -v grep | awk '{print $2}' | xargs kill; nohup python3.11 ~/code/volume-control/src/volume_control.py &> /tmp/nohup.out & disown
